# examples/meme-service/server.service.in
[Unit]
Description=%NAME% service
After=syslog.target network.target

[Service]
ExecStart=/opt/%NAME%/%NAME%-service
WorkingDirectory=/opt/%NAME%
User=%NAME%
Restart=always
RestartSec=0

[Install]
WantedBy=multi-user.target
